# supervisor/phpapp.conf
[program:phpapp]
command=docker run -p 80:80 phpapp
autostart=true
autorestart=true
stderr_logfile=/opt/run/webphp/log/phpapp.err.log
stdout_logfile=/opt/run/webphp/log/phpapp.out.log
